# randombg-timer/.config/systemd/user/randombg.service
[Unit]
Description=Set a random Sway background from Nextcloud

[Service]
ExecStart=/usr/local/bin/randombg
Type=oneshot
